# nexus.service
[Unit]
Description=NEXUS Multi-Channel Agentic Service
After=network.target

[Service]
Type=simple
User=puretensorai
Group=puretensorai
WorkingDirectory=/home/puretensorai/nexus
ExecStart=/usr/bin/python3 /home/puretensorai/nexus/nexus.py
Restart=always
RestartSec=10
Environment=HOME=/home/puretensorai
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
